# scripts/b70/b70-watchdog.service
# b70-watchdog.service — single watchdog pass (driven by b70-watchdog.timer).
#
# Detects B70 drops at runtime, captures a diagnostic snapshot, stops vLLM,
# attempts non-destructive recovery only, and alerts. Never re-binds silently
# under a running vLLM and never does a PCI reset/bridge SBR.

[Unit]
Description=B70 runtime health watchdog (one pass)
After=intel-xe-power-cap.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/b70-watchdog
# Optional: point at an alerting script (Slack/ntfy/email). Called as:
#   <hook> <verdict> <snapshot-path>
# Environment=B70_ALERT_HOOK=/usr/local/sbin/b70-alert
Environment=B70_VLLM_UNIT=b70-vllm.service
Environment=B70_STATE_DIR=/var/log/b70
