# scripts/b70/b70-watchdog.timer
# b70-watchdog.timer — run the B70 watchdog every 30s.
#
# Install via ./install.sh, then:
#   sudo systemctl enable --now b70-watchdog.timer

[Unit]
Description=Run B70 health watchdog periodically

[Timer]
OnBootSec=60
OnUnitActiveSec=30
AccuracySec=5
Unit=b70-watchdog.service

[Install]
WantedBy=timers.target
